# APSetup/Files/OrangePi/etc/update-motd.d/01-phoebus-header
#!/bin/bash
#
# Copyright (c) Authors: https://www.armbian.com/authors
#
# This file is licensed under the terms of the GNU General Public
# License version 2. This program is licensed "as is" without any
# warranty of any kind, whether express or implied.

# DO NOT EDIT THIS FILE but add config options to /etc/default/armbian-motd
# any changes will be lost on board support package update

THIS_SCRIPT="header"
MOTD_DISABLE=""

VERSION=1.0.0

for f in $MOTD_DISABLE; do
        [[ $f == $THIS_SCRIPT ]] && exit 0
done


KERNELID=$(uname -r)

TERM=linux toilet -f standard -F metal $(echo "Phoebus")
echo -e "Welcome to \e[0;91mPhoebus ${VERSION} Linux $KERNELID\x1B[0m\n"